# deployment/systemd/nexus-pisowifi.service
[Unit]
Description=NEXUS PISOWIFI System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/NEXUS-PISOWIFI
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=CAPTIVE_FORCE_REDIRECT=true
Environment=CAPTIVE_SKIP_AP_VERIFY=true
Environment=CAPTIVE_COUNTRY_CODE=PH
Environment=PATH=/usr/local/bin:/usr/bin
ExecStart=/usr/bin/npm run start:portal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
